# providers/os/resources/services/testdata/systemd/etc/systemd/system/explicit-socket.socket
[Unit]
Description=Explicit Socket

[Socket]
ListenStream=/run/explicit/socket
Service=explicit-socket-service.service
