# providers/os/resources/services/testdata/systemd/etc/systemd/system/intermediate-dep-require.service
[Unit]
Description=Intermediate Dep Service
