# lumi/resources/services/testdata/systemd/usr/lib/systemd/system/default.target
[Unit]
Description=Default Target
Wants=aliastest.service intermediate.target implicit-socket.socket explicit-socket.socket
