# lumi/resources/services/testdata/systemd/usr/lib/systemd/system/intermediate.target
[Unit]
Description=Intermediate Target
Wants=masked.service
