# lumi/resources/services/testdata/systemd/usr/lib/systemd/system/masked.service
[Unit]
Description=Masked service
Wants=not-enabled.service
