# resources/packs/os/services/testdata/systemd/etc/systemd/system/aliased.service
[Unit]
Description=Aliased Service
Wants=aliased-wants.service aliased-missing.service
Requires=aliased-requires.service
