# resources/packs/os/services/testdata/systemd/etc/systemd/system/implicit-socket.socket
[Unit]
Description=Implicit Socket

[Socket]
ListenStream=/run/implicit/socket
