# lumi/resources/services/testdata/systemd/etc/systemd/system/aliased-wants.service
[Unit]
Description=Aliased Wants Service
