# lumi/resources/services/testdata/systemd/etc/systemd/system/explicit-socket-service.service
[Unit]
Description=Explicit Socket Service
