# lumi/resources/services/testdata/systemd/etc/systemd/system/implicit-socket.service
[Unit]
Description=Implicit Socket Service
